# monitoring/cgroups/simulator.slice
[Unit]
Description=Simulator workload (highest priority, no CPU limit)

[Slice]
CPUWeight=1000
