# monitoring/cgroups/streamer.slice
[Unit]
Description=Streamer workload (max 60% CPU)

[Slice]
CPUWeight=300
CPUQuota=960%
