# monitoring/cgroups/tools.slice
[Unit]
Description=Tools workload (max 20% CPU)

[Slice]
CPUWeight=100
CPUQuota=320%
